# rootfs/etc/systemd/system/mhi-qmi@.service
[Unit]
Description=MHI QMI port %I
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=simple
ExecStart=/bin/cat /dev/%I
