# rootfs/lib/systemd/system/at_server.service
[Unit]
Description=at_server

[Service]
Type=simple
ExecStart=/usr/bin/at-server
Restart=always

[Install]
WantedBy=multi-user.target
