# rootfs/lib/systemd/system/diag-relay.service
[Unit]
Description=MDM Diag Relay Service

[Service]
Type=simple
ExecStart=/usr/bin/diag-relay
Restart=always

[Install]
WantedBy=multi-user.target
